# monitoring/systemd/vps-health.timer
[Unit]
Description=Run VPS health check

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
