# monitoring/systemd/vps-health.service
[Unit]
Description=VPS health check

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vps-health
